# camera_data_extractor/rootfs/etc/services.d/camera-extractor/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Take down the S6 supervision tree when the service exits
# ==============================================================================

declare exit_code=${1}

if [[ "${exit_code}" -ne 0 ]] && [[ "${exit_code}" -ne 256 ]]; then
    bashio::log.error "Camera Data Extractor exited with code ${exit_code}"
    exec /run/s6/basedir/bin/halt
fi
